# data/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
nodaemon=true

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx

[program:php-fpm]
command=/usr/sbin/php-fpm